# systemd/things.target
[Unit]
Wants=things-web.target things-batch.target things-interactive.target things-immediate.target

[Install]
WantedBy=multi-user.target
